# jupiter-fan-control.service
[Unit]
Description=Jupiter fan control

[Service]
Type=simple
#StandardOutput=null
#KillSignal=SIGINT
ExecStart=/usr/share/jupiter-fan-control/fancontrol.py

# disable buffering so STDOUT and STDERR show up immediately in logs
Environment=PYTHONUNBUFFERED=1

OOMScoreAdjust=-1000
Restart=on-failure

[Install]
WantedBy=multi-user.target
